# rootfs/etc/systemd/system/heartbeat.service
[Unit]
Description=ICPC heartbeat

[Service]
Type=oneshot
ExecStart=/opt/heartbeat