# rootfs/etc/systemd/system/heartbeat.service.d/override.conf
[Service]
Environment="HEARTBEATURL=http://192.168.99.2:3000/"